# disk_image/create/jetson_orin_nano/disk_template/APP/etc/systemd/system/dt_init.service
[Unit]
Description=This script setups the robot on boot
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /usr/local/bin/dt-init

[Install]
WantedBy=multi-user.target